# config/systemd/user/btrfs-snp.service.example
[Unit]
Description="btrfs-snp-wrapper service unit file"

[Service]
ExecStart="/home/deck/Applications/btrfs-snp/bin/btrfs-snp-wrapper"
