# config/systemd/user/btrfs-snp.timer
[Unit]
Description="Run btrfs-snp (wrapper) script hourly"

[Timer]
OnCalendar=*-*-* *:00:00
Unit=btrfs-snp.service

[Install]
WantedBy=multi-user.target
